# socks5/rootfs/etc/services.d/example/run
#!/bin/sh
set -e

echo "Starting SOCKS5 Proxy Server..."
exec /usr/bin/docker run --rm -p 1080:1080 \
  -e PROXY_USER=${PROXY_USER} \
  -e PROXY_PASSWORD=${PROXY_PASSWORD} \
  serjs/go-socks5-proxy
